# supabase/gh-runner@.service
# Systemd template unit for an ephemeral GitHub Actions runner (vanilla actions/runner).
#
# This is a TEMPLATE: install.sh sed-substitutes __RUNNER_BASE__ and __RUN_USER__ and writes the
# result to /etc/systemd/system/gh-runner-supabase@.service (per-type name — avoids collision when
# light and supabase runners share a host), then enables one instance per runner dir
# (gh-runner-supabase@1, gh-runner-supabase@2, ...). Each instance %i maps to __RUNNER_BASE__/%i,
# which holds its own actions/runner install + config.env (mode 600) + runner-loop.sh.
#
# Do not install this template by hand; run install.sh.

[Unit]
Description=Ephemeral GitHub Actions runner (instance %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__RUNNER_BASE__/%i
# The runner derives RUNNER_TOOL_CACHE from AGENT_TOOLSDIRECTORY, so setup-python/-node/-deno resolve
# tools from a shared, persistent cache (hosted-runner parity). REQUIRED for setup-python on
# non-Ubuntu hosts, which can't download a prebuilt Python. install.sh substitutes the path.
Environment=AGENT_TOOLSDIRECTORY=__TOOLCACHE_DIR__
ExecStart=/bin/bash __RUNNER_BASE__/%i/runner-loop.sh
Restart=always
RestartSec=10
# Give runner-loop.sh time to deregister (SIGTERM trap) before systemd SIGKILLs it.
KillSignal=SIGTERM
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
