# light/gh-runner@.service
# Systemd template unit for an ephemeral GitHub Actions runner (vanilla actions/runner).
#
# This is a TEMPLATE: install.sh sed-substitutes __RUNNER_BASE__, __RUN_USER__, __TOOLCACHE_DIR__,
# and __PLAYWRIGHT_ENV__ (empty unless --with-playwright) and writes the result to
# /etc/systemd/system/gh-runner-light@.service
# (per-type name — avoids collision when light and supabase runners share a host), then enables
# one instance per runner dir (gh-runner-light@1, gh-runner-light@2, ...). Each instance %i maps
# to __RUNNER_BASE__/%i, which holds its own actions/runner install + config.env + fleet scripts.
#
# ExecStart targets runner-bootstrap.sh (NOT runner-loop.sh directly). The bootstrap is the
# stable supervisor entrypoint that manages the crash counter and last_good/ rollback; it exec's
# runner-loop.sh as the swappable payload. Fleet-code updates only touch the payload.
#
# Do not install this template by hand; run install.sh.

[Unit]
Description=Ephemeral GitHub Actions runner (instance %i)
After=network-online.target
Wants=network-online.target
# WHY: disable the start-rate-limiter so systemd cannot mark the unit failed and stop relaunches
# before runner-bootstrap.sh's crash counter reaches ROLLBACK_THRESHOLD and triggers last_good/
# restore. RestartSec=10 in [Service] still throttles rapid restarts.
StartLimitIntervalSec=0

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__RUNNER_BASE__/%i
# The runner derives RUNNER_TOOL_CACHE from AGENT_TOOLSDIRECTORY, so setup-python/-node/-deno resolve
# tools from a shared, persistent cache (hosted-runner parity). REQUIRED for setup-python on
# non-Ubuntu hosts, which can't download a prebuilt Python. install.sh substitutes the path.
Environment=AGENT_TOOLSDIRECTORY=__TOOLCACHE_DIR__
# CRITICAL — per-instance HOME. Co-hosted runners run as the SAME user, so a single shared HOME means
# concurrent jobs race on per-user caches: pnpm/action-setup (/setup-pnpm), setup-node's registry-auth
# .npmrc, and the deno/pip/npm caches. That race already corrupted pnpm's install (empty dist/worker.js)
# fleet-wide. Giving each instance %i its own HOME isolates those caches. install.sh creates the dir.
# (The shared toolcache + Playwright browser cache stay shared on purpose — set via their own env.)
Environment=HOME=__RUNNER_BASE__/%i/home
# Playwright capability (only when installed with --with-playwright): points every job at a shared,
# persistent browser cache so `playwright install` is a cache hit and needs no root. Blank otherwise.
__PLAYWRIGHT_ENV__
ExecStart=/bin/bash __RUNNER_BASE__/%i/runner-bootstrap.sh
Restart=always
RestartSec=10
# Give runner-loop.sh time to deregister (SIGTERM trap) before systemd SIGKILLs it.
KillSignal=SIGTERM
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
