# ios/com.example.gh-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  com.example.gh-runner.plist — LaunchAgent for the iOS/Android-on-Mac self-hosted runner.

  IMPORTANT: This file is a TEMPLATE. install.sh performs three sed substitutions before
  writing to ~/Library/LaunchAgents/:
    __SERVICE_LABEL__ -> the LaunchAgent label (default: com.example.gh-runner; set via
                         --service-label or SERVICE_LABEL env at install time)
    __RUNNER_DIR__    -> absolute path of the runner install dir (default: ~/actions-runner-e2e,
                         expanded to the real path by install.sh — e.g. /Users/alice/actions-runner-e2e)
    __LOG_DIR__       -> absolute path of the log dir (default: ~/Library/Logs, expanded likewise)

  Do not load this template directly; always go through install.sh.

  NOTE: pass --service-label com.yourteam.gh-runner to install.sh to customize the label
  (e.g. com.acme.gh-runner). The template filename is always com.example.gh-runner.plist.

  ProgramArguments targets runner-bootstrap.sh (NOT runner-loop.sh directly). The bootstrap is
  the stable supervisor entrypoint that manages the crash counter and last_good/ rollback; it
  exec's runner-loop.sh as the swappable payload. Fleet-code updates only touch the payload.

  KeepAlive=true (bare boolean) causes launchd to restart runner-bootstrap.sh on ANY exit,
  including clean exit 0 (which is how an applied fleet-code update signals the supervisor to
  relaunch with the new payload). IMPORTANT: a SuccessfulExit=false dict would NOT relaunch
  on clean exit and would break the update mechanism — do not change this form.
-->
<dict>
  <key>Label</key>
  <string>__SERVICE_LABEL__</string>

  <key>ProgramArguments</key>
  <array>
    <!-- WHY: invoke bash explicitly so launchd doesn't require runner-bootstrap.sh to be in $PATH.
         LaunchAgent environments have a minimal PATH and may not find the shebang interpreter. -->
    <string>/bin/bash</string>
    <string>__RUNNER_DIR__/runner-bootstrap.sh</string>
  </array>

  <!-- Start immediately when the plist is loaded and on every user login. -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart runner-loop.sh automatically if it exits unexpectedly. -->
  <key>KeepAlive</key>
  <true/>

  <!-- Throttle restarts: if runner-loop.sh exits immediately (e.g. config.env missing), don't
       respawn faster than every 30s. Prevents a tight crash-loop from flooding the logs. -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <!-- WorkingDirectory ensures relative paths inside runner-loop.sh resolve correctly.
       NOTE: This is also rewritten from __RUNNER_DIR__ by install.sh. -->
  <key>WorkingDirectory</key>
  <string>__RUNNER_DIR__</string>

  <!-- NOTE: stdout and stderr go to separate files for unambiguous log tailing.
       'tail -f ~/Library/Logs/gh-runner.out.log' to follow runner output. -->
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/gh-runner.out.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/gh-runner.err.log</string>

</dict>
</plist>
